# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:nodejs]
command=cd /app/pdf-node && yarn start

; [program:java]
; command=java -jar /app.jar
